# mangabot.service
[Unit]
Description=MangaBot - JMComic下载机器人
Documentation=https://github.com/your-repo/mangabot
After=network.target
Wants=network.target

[Service]
Type=simple
User=mangabot
Group=mangabot
WorkingDirectory=/opt/mangabot
ExecStart=/usr/bin/python3 /opt/mangabot/bot.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

# 环境变量配置
Environment=PYTHONUNBUFFERED=1
Environment=MANGA_DOWNLOAD_PATH=/var/lib/mangabot/downloads
Environment=NAPCAT_WS_URL=ws://localhost:8080/qq

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/mangabot/downloads /var/log/mangabot

# 资源限制
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target